# os/etc/systemd/system/falco-modern-bpf.service.d/10-atomic.conf
# Atomic OS drop-in for falco-modern-bpf.service
#
# Overrides from the upstream Falco RPM unit:
#   1. ExecStart: pass --modern-bpf explicitly so Falco cannot silently
#      fall back to the classic eBPF driver (which requires a DKMS-built
#      .o probe file that does not exist on an immutable OS).
#      -c /etc/falco/falco.yaml ensures our config (engine, rules, webhook)
#      is used instead of the RPM default.
#   2. StandardOutput/StandardError: upstream unit sets StandardOutput=null
#      which also silences StandardError via inheritance.  Restore journal
#      logging so Falco alerts are visible via journalctl and picked up by
#      the system log (needed for demo and for production audit trail).

[Service]
# Empty assignment clears the inherited ExecStart, then we set our own.
ExecStart=
ExecStart=/usr/bin/falco --modern-bpf -c /etc/falco/falco.yaml

# Restore journal logging.
StandardOutput=journal
StandardError=journal
