# os/etc/systemd/system/falco-bpf.service.d/10-atomic.conf
# Atomic OS drop-in for falco-bpf.service
#
# Overrides from the upstream Falco RPM unit:
#   1. ExecStart: use our config explicitly (-c /etc/falco/falco.yaml).
#      The probe path is set in the config file (no FALCO_BPF_PROBE env var
#      needed — probe is pre-compiled into the image at build time).
#   2. StandardOutput/StandardError: upstream unit sets StandardOutput=null
#      which also silences StandardError.  Restore journal logging so Falco
#      alerts and errors are visible via journalctl.
#
# NOTE: No ExecStartPre for falcoctl driver install — the eBPF probe is
# compiled into the image at /usr/share/atomic/falco/falco-bpf.o by the
# agentic-os Containerfile.  Runtime probe download is not needed.

[Service]
# Override ExecStart to use our config explicitly.
# Empty assignment first clears the inherited value.
ExecStart=
ExecStart=/usr/bin/falco -c /etc/falco/falco.yaml

# Restore journal logging (upstream sets StandardOutput=null which
# also silences StandardError via inheritance).
StandardOutput=journal
StandardError=journal
