# os/etc/systemd/system/falco-bpf.service.d/10-atomic.conf
# Atomic OS drop-in for falco-bpf.service
#
# Overrides from the upstream Falco RPM unit:
#   1. ExecStartPre: download the kernel-specific eBPF probe via falcoctl
#      before Falco tries to load it.  This handles first-boot and kernel
#      upgrades automatically.
#   2. ExecStart override: add -c to use our config (explicit > implicit).
#   3. StandardOutput/StandardError: upstream unit discards all output
#      (StandardOutput=null).  We restore journal logging so Falco alerts
#      and errors are visible via journalctl.

[Service]
# Download (or update) the eBPF probe for the currently running kernel.
# falcoctl writes the probe to $HOME/.falco/<version>/<arch>/<krel>/falco.o
# and Falco picks it up automatically when engine.kind=ebpf is set.
ExecStartPre=/usr/bin/falcoctl driver install --type ebpf

# Use our config explicitly (default is also /etc/falco/falco.yaml, but
# being explicit makes overrides and debugging easier).
ExecStart=
ExecStart=/usr/bin/falco -c /etc/falco/falco.yaml -o engine.kind=ebpf

# Restore journal logging (upstream sets StandardOutput=null which also
# silences StandardError, hiding all Falco output from journalctl).
StandardOutput=journal
StandardError=journal
